# selection-panel.service
[Unit]
Description=Auswahlpanel v2.2.0 (Serial→WebSocket)
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/selection-panel
ExecStart=/home/pi/selection-panel/venv/bin/python server.py
Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal

# Security Hardening
NoNewPrivileges=yes
ProtectSystem=strict
PrivateTmp=yes
ReadWritePaths=/home/pi/selection-panel

[Install]
WantedBy=multi-user.target
